# Subir a otro Servidor/systemd/redlibre-uexpired.service
[Unit]
Description=RedLibre — cortar/contar usuarios SSH vencidos (vence + chage)

[Service]
Type=oneshot
ExecStart=/bin/uexpired
StandardOutput=null
StandardError=journal
