# Subir a otro Servidor/systemd/redlibre-uexpired.service
[Unit]
Description=RedLibre — actualizar contador de usuarios expirados (/etc/RedLibre/Exp)

[Service]
Type=oneshot
ExecStart=/bin/uexpired
StandardOutput=null
StandardError=journal
